# lib/bin/standby.crontab
5 * * * *  /bin/bash /shared/bin/backup-db.sh testdb
