# lib/bin/standby.crontab
20 * * * *  /bin/bash /vagrant_lib/bin/backup-db.sh testdb
